# kafka-jdk6/kafka.0.8.2.2/resources/kafka/kafka.conf
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autorestart=true
stdout_logfile=/var/log/kafka.log
stderr_logfile=/var/log/kafka.err
